# scripts/systemd/tactical-frontend.service
[Unit]
Description=Tactical Command Frontend Service
After=network.target tactical-backend.service
Wants=tactical-backend.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/app/frontend

# Environment
Environment="HOST=0.0.0.0"
Environment="PORT=3000"

# Start command
ExecStart=/usr/bin/yarn start

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tactical-frontend

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target