# scripts/systemd/tactical-backend.service
[Unit]
Description=Tactical Command Backend Service
After=network.target mongodb.service
Wants=mongodb.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/app/backend

# Virtual environment activation and start
ExecStart=/root/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8001 --workers 1

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tactical-backend

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target